# php_apache/root/etc/services.d/httpd/run
#!/usr/bin/with-contenv sh
exec httpd -DFOREGROUND

